# vpn-watchdog.service
[Unit]
Description=VPN Watchdog Service
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/vpn-watchdog.sh
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vpn-watchdog

# Security settings
User=root
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/v2ray /opt/outline /var/log

[Install]
WantedBy=multi-user.target